# deploy/local/com.jobradar.emailagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.jobradar.emailagent</string>

  <!-- the pipx-installed console script; replace the path if pipx installs elsewhere -->
  <key>ProgramArguments</key>
  <array>
    <string>/Users/duaneo/.local/bin/job-radar-agent</string>
    <string>run</string>
    <string>--once</string>
  </array>

  <!-- where .env + data/ live (durable, survives pipx upgrade) -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>AGENT_HOME</key>
    <string>/Users/duaneo/Library/Application Support/JobRadarAgent</string>
  </dict>

  <key>StartInterval</key>
  <integer>900</integer>            <!-- every 15 min -->
  <key>RunAtLoad</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Users/duaneo/Library/Logs/jobradar-emailagent.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/duaneo/Library/Logs/jobradar-emailagent.log</string>
</dict>
</plist>
